# ton-launch-tracker.service
[Unit]
Description=TON DeDust + STON.fi New Launch Telegram Tracker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/ton-launch-tracker
EnvironmentFile=/root/ton-launch-tracker/.env
ExecStart=/root/ton-launch-tracker/.venv/bin/python -m tracker
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
